# vault-sup.conf
[program:vault]
command=vault server -config=/etc/vault.conf
user=root
